# keygen/configs/cron.daily
#! /bin/sh

# Do a periodic trustdb maintenance.  Use --batch, so the maintenance is done
# only when needed.
runuser -u copr-signer -- /usr/bin/gpg-copr --check-trustdb --batch
