# backend/conf/crontab/copr-backend
#!/usr/bin/sh

runuser -c "PYTHONUNBUFFERED=1 python /usr/bin/copr_prune_results.py &> /dev/null" - copr
