# backend/conf/crontab/copr-backend
#!/usr/bin/sh

runuser -c "PYTHONUNBUFFERED=1 python /usr/share/copr/copr_prune_results.py" - copr >> /var/log/copr/prune_old.log 2>&1
