# backend/conf/crontab/weekly
#! /bin/sh

runuser -c "ionice --class idle /usr/bin/copr-backend-analyze-results && /usr/bin/copr-backend-generate-graphs" - copr >&2 >/dev/null
