# backend/conf/crontab/weekly
#! /bin/sh

runuser -c "/usr/bin/copr-backend-analyze-results && /usr/bin/copr-backend-generate-graphs" - copr >&2 >/dev/null
