# frontend/conf/cron.daily/copr-frontend
#!/usr/bin/sh

runuser -c '/usr/share/copr/coprs_frontend/manage.py vacuum_graphs' - copr-fe
