# frontend/conf/cron.daily/copr-frontend
#!/usr/bin/sh

runuser -c '/usr/share/copr/coprs_frontend/run/update_graphs_data.py --remove' - copr-fe
