# backend/crontab/copr-backend
#!/usr/bin/sh

/usr/sbin/copr-prune-repo >/var/log/copr.log 2>&1
