# backend/crontab/copr-backend
#!/usr/bin/sh

python /usr/share/copr/copr_prune_results.py >> /var/log/copr/prune_old.log 2>&1
